# agent/systemd/scrumban-agent.service
[Unit]
Description=Scrumban Agent V2 (cliente Claude Code remoto)
Documentation=https://github.com/devari/scrumban-backend-v2
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=scrumban-agent
Group=scrumban-agent
ExecStart=/usr/bin/node /opt/scrumban-agent/dist/index.js
Restart=always
RestartSec=5
WorkingDirectory=/var/lib/scrumban-agent
Environment=NODE_ENV=production
# Credenciais do Claude Code CLI (ANTHROPIC_API_KEY etc.). O `-` prefix
# faz o systemd NÃO falhar se o arquivo ainda não existir — install.sh
# cria um placeholder e o operador preenche manualmente.
EnvironmentFile=-/etc/scrumban-agent/environment
StandardOutput=journal
StandardError=journal

# Hardening (systemd security)
#
# NoNewPrivileges=false é necessário para o handler SET_ENV (Fase 4
# task plan-2026-05-13) conseguir rodar `sudo systemctl restart
# scrumban-agent` via sudoers entry de escopo mínimo. Com `true`, o
# systemd bloqueia QUALQUER escalação — inclusive a do sudoers, que é
# segura (rule única, sem wildcards). Trade-off:
#   - Risco: se o processo do agent for comprometido, ele PODE chamar
#     `sudo` — mas o sudoers só permite EXATAMENTE
#     `/usr/bin/systemctl restart scrumban-agent`. Blast radius máximo
#     = DoS por loop de restart no próprio serviço. Sem acesso root.
#   - Sem isso: SET_ENV nunca aplica novas credenciais sem intervenção
#     manual (defeito de produto).
NoNewPrivileges=false
PrivateTmp=true
ProtectSystem=strict
# ProtectHome=no permite que o agente (e o claude filho) escrevam em
# /home/dev/projetos (workspace dos projetos gerenciados). Com read-only
# o systemd bloqueia escrita mesmo que ReadWritePaths liste o path.
# Leitura de ~/.claude/CLAUDE.md continua funcionando normalmente.
ProtectHome=no
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
LockPersonality=true
RestrictRealtime=true
RestrictNamespaces=true
# Apenas paths estritamente necessários ficam read-write
ReadWritePaths=/etc/scrumban-agent /var/lib/scrumban-agent /var/log/scrumban-agent /home/dev/projetos

# Limites de recursos
MemoryMax=512M
TasksMax=100

[Install]
WantedBy=multi-user.target
